# conf/supervisord-ws-server.conf
[program:ws-server]
command=php /app/ws-socket.php
numprocs=1
autostart=true
autorestart=true
